# jmeter-icap/scripts/flask_scaled.service
# /etc/systemd/system/flask_scaled.service
[Unit]
Description=WSGI App for ICAP Testing UI Front End
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/glasswall/p-k8-jmeter-test-engine/jmeter-icap/scripts
ExecStart=/home/glasswall/p-k8-jmeter-test-engine/jmeter-icap/scripts/exec_scaled.sh
Restart=always

[Install]
WantedBy=multi-user.target
